# backporcher-dashboard.service.example
# backporcher-dashboard.service — Standalone dashboard web server
# Runs the aiohttp dashboard as a separate service (e.g., behind a reverse proxy).
#
# Replace YOUR_USER, YOUR_GROUP, and paths below to match your deployment.
#
# Install:
#   sudo cp backporcher-dashboard.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now backporcher-dashboard

[Unit]
Description=Backporcher Dashboard
After=network.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_GROUP
WorkingDirectory=/home/YOUR_USER/backporcher

ExecStart=/home/YOUR_USER/backporcher/scripts/start-dashboard.sh

Restart=on-failure
RestartSec=5

# Environment — edit these for your setup
Environment=BACKPORCHER_BASE_DIR=/home/YOUR_USER/backporcher
Environment=BACKPORCHER_GITHUB_OWNER=your-github-org
Environment=BACKPORCHER_DASHBOARD_PORT=8085
Environment=BACKPORCHER_DASHBOARD_HOST=127.0.0.1
Environment=PATH=/home/YOUR_USER/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Load password from a credential file. Adjust path to your credential store.
# LoadCredential=backporcher-dashboard-password:/path/to/credentials/backporcher-dashboard-password
# Or set directly (less secure):
# Environment=BACKPORCHER_DASHBOARD_PASSWORD=changeme

# Security hardening
ProtectSystem=full
ProtectHome=no
ReadWritePaths=/home/YOUR_USER/backporcher/data
PrivateTmp=yes
PrivateDevices=yes
PrivateIPC=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
RestrictNamespaces=yes
NoNewPrivileges=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=backporcher-dashboard

[Install]
WantedBy=multi-user.target
